# keboola-config/supervisord/services/_continuous-pull.conf
; Wrapper-owned processes for continuous-pull. The watched app's processes
; are defined in its own keboola-config/supervisord/services/*.conf files
; (any number of [program:*] sections -- all are loaded side by side).
;
; Naming note: the wrapper's setup.sh copies this file into the user's
; supervisord services dir verbatim. The leading underscore keeps it sorted
; first and avoids colliding with the conventional user-owned "app.conf".

[program:pull-loop]
command=bash /tmp/continuous-pull/scripts/pull_loop.sh
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pull-api]
command=python3 /tmp/continuous-pull/scripts/api_server.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
